# environments/dev/adakings-backend-dev.service
[Unit]
Description=Adakings Backend API - Development Environment
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/adakings/dev
Environment=DJANGO_ENVIRONMENT=dev
EnvironmentFile=/var/www/adakings/dev/.env
ExecStart=/var/www/adakings/dev/venv/bin/gunicorn \
    --config /var/www/adakings/dev/environments/dev/gunicorn.conf.py \
    adakings_backend.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5

# Security settings (relaxed for development)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/adakings /var/log/adakings /tmp
PrivateDevices=false  # Allow access to devices for development
ProtectKernelTunables=false  # Relaxed for development
ProtectKernelModules=false  # Relaxed for development
ProtectControlGroups=false  # Relaxed for development

[Install]
WantedBy=multi-user.target
